# momoclo/supervisord.conf
[program:node]
environment = NODE_ENV=production
command = node server.js
directory = /home/dotcloud/current
